# deploy/home-dashboard.service.in
[Unit]
Description=Home Dashboard service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Group=@USER@
# lgpio creates its .lgd-nfy notification pipe in the current directory.
# Keep the runtime in the writable data path rather than the read-only app tree.
WorkingDirectory=/var/lib/home-dashboard
Environment=HOME_DASHBOARD_DATA=/var/lib/home-dashboard
Environment=HOME_DASHBOARD_STATIC=/opt/home-dashboard/frontend/dist
Environment=HOME_DASHBOARD_HOST=0.0.0.0
Environment=HOME_DASHBOARD_PORT=8765
Environment=PYTHONPATH=/opt/home-dashboard/backend
Environment=PYTHONUNBUFFERED=1
Environment=XDG_RUNTIME_DIR=/run/user/@UID@
Environment=PIPEWIRE_RUNTIME_DIR=/run/user/@UID@
Environment=PULSE_RUNTIME_PATH=/run/user/@UID@/pulse
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/@UID@/bus
# Raspberry Pi OS labwc normally uses wayland-0; the app also discovers any
# active wayland-* socket at runtime for Wayfire and future compositor changes.
Environment=WAYLAND_DISPLAY=wayland-0
Environment=GPIOZERO_PIN_FACTORY=lgpio
ExecStart=/opt/home-dashboard/.venv/bin/home-dashboard
Restart=always
RestartSec=3
TimeoutStopSec=20
# The service uses two sudoers-allowlisted helpers for Wi-Fi and self-restart.
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/home-dashboard @HOME@/.local/bin

[Install]
WantedBy=multi-user.target
