# deploy/home-dashboard-update.service
[Unit]
Description=Update BodaHomeDash from GitHub main
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/lib/home-dashboard/update.sh
Nice=10
